# start-api.service
[Unit]
Description=Trading ML API Service
Documentation=https://github.com/yourusername/trading-ml-api
After=network.target

[Service]
Type=notify
User=www-data
WorkingDirectory=/opt/trading-api
ExecStart=/opt/trading-api/venv/bin/python -m uvicorn api.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --access-log

ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-api

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/trading-api/data /opt/trading-api/logs

# Resource limits
LimitNOFILE=65535
LimitNPROC=512

[Install]
WantedBy=multi-user.target
